# .config/systemd/user/duplicity.service
[Unit]
Description=Run duplicity to back up my home directory

[Service]
ExecStart=/home/spwhitton/bin/duply-run
