# .config/systemd/user/i3-rotate-wallpaper.service
[Unit]
Description=Rotate i3 wallpaper

[Service]
Type=oneshot
ExecStart=%h/bin/i3-rotate-wallpaper
