# .config/systemd/user/i3-rotate-wallpaper.timer
[Unit]
Description=Run i3-rotate-wallpaper daily

[Timer]
OnActiveSec=1 day
OnUnitActiveSec=1 day
